# supervisord.conf
; ═══════════════════════════════════════════════════════════════════════
; MIKKIE WORLD — Supervisord Configuratie
; Industriestandaard process manager — geen os.fork() nodig
;
; Installeer op macOS: brew install supervisor
; (NIET pip3 install — macOS blokkeert dat via PEP 668)
;
; Start:      supervisord -c ~/mikkieworld/supervisord.conf
; Status:     supervisorctl -c ~/mikkieworld/supervisord.conf status
; Stop:       supervisorctl -c ~/mikkieworld/supervisord.conf stop all
;             supervisorctl -c ~/mikkieworld/supervisord.conf shutdown
; ═══════════════════════════════════════════════════════════════════════

[supervisord]
logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=%(ENV_HOME)s/mikkieworld/pids/supervisord.pid
nodaemon=false
silent=false
user=%(ENV_USER)s

; ─── Web interface (optioneel) ─────────────────────────────────────────
[inet_http_server]
port=127.0.0.1:9001
; username=mikkie
; password=geheim

[supervisorctl]
serverurl=http://127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ─── BRAIN — Centrale Orchestrator ────────────────────────────────────
[program:brain]
command=python3 %(ENV_HOME)s/mikkieworld/mikkie_brain.py loop
directory=%(ENV_HOME)s/mikkieworld
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stderr_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/brain_stderr.log
stdout_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/brain_stdout.log
stderr_logfile_maxbytes=5MB
stdout_logfile_maxbytes=5MB
environment=XAI_API_KEY="%(ENV_XAI_API_KEY)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",TELEGRAM_CHAT_ID="%(ENV_TELEGRAM_CHAT_ID)s",GUMROAD_API_TOKEN="%(ENV_GUMROAD_API_TOKEN)s"

; ─── TELEGRAM COMMANDER — Mobiele Controle ────────────────────────────
[program:commander]
command=python3 %(ENV_HOME)s/mikkieworld/mikkie_telegram_commander.py
directory=%(ENV_HOME)s/mikkieworld
autostart=true
autorestart=true
startretries=5
startsecs=3
stopwaitsecs=10
stderr_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/commander_stderr.log
stdout_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/commander_stdout.log
stderr_logfile_maxbytes=5MB
stdout_logfile_maxbytes=5MB
environment=XAI_API_KEY="%(ENV_XAI_API_KEY)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",TELEGRAM_CHAT_ID="%(ENV_TELEGRAM_CHAT_ID)s"

; ─── HEALER — Zelf-herstellend Systeem ────────────────────────────────
[program:healer]
command=python3 %(ENV_HOME)s/mikkieworld/mikkie_healer.py start
directory=%(ENV_HOME)s/mikkieworld
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=30
stderr_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/healer_stderr.log
stdout_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/healer_stdout.log
stderr_logfile_maxbytes=5MB
stdout_logfile_maxbytes=5MB
environment=XAI_API_KEY="%(ENV_XAI_API_KEY)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",TELEGRAM_CHAT_ID="%(ENV_TELEGRAM_CHAT_ID)s"

; ─── Groepsbeheer ──────────────────────────────────────────────────────
[group:mikkie_daemons]
programs=brain,commander,healer
priority=999

[program:agent]
command=/opt/homebrew/bin/python3 %(here)s/mikkie_agent.py
directory=%(here)s
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=%(here)s/agent.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=%(here)s/agent_error.log

[program:mikkie_social_poster]
command=python3 /Users/hendrikbroeze/mikkie/scripts/social_poster.py
directory=/Users/hendrikbroeze/mikkie
autostart=true
autorestart=true
startretries=5
stderr_logfile=/Users/hendrikbroeze/mikkie/logs/social_poster_err.log
stdout_logfile=/Users/hendrikbroeze/mikkie/logs/social_poster_out.log

[program:mikkie_analytics]
command=python3 /Users/hendrikbroeze/mikkie/scripts/analytics_collector.py
directory=/Users/hendrikbroeze/mikkie
autostart=true
autorestart=true
startretries=5
stderr_logfile=/Users/hendrikbroeze/mikkie/logs/analytics_err.log
stdout_logfile=/Users/hendrikbroeze/mikkie/logs/analytics_out.log

[program:mikkie_social_healer]
command=python3 /Users/hendrikbroeze/mikkie/scripts/self_healer.py
directory=/Users/hendrikbroeze/mikkie
autostart=true
autorestart=true
startretries=10
stderr_logfile=/Users/hendrikbroeze/mikkie/logs/healer_err.log
stdout_logfile=/Users/hendrikbroeze/mikkie/logs/healer_out.log
