# supervisord.conf
; ═══════════════════════════════════════════════════════════════════════
; MIKKIE WORLD — Supervisord Configuratie
; Industriestandaard process manager — geen os.fork() nodig
;
; Installeer: pip3 install supervisor
; Start:      supervisord -c ~/mikkieworld/supervisord.conf
; Status:     supervisorctl status
; Stop:       supervisorctl stop all && supervisorctl shutdown
; ═══════════════════════════════════════════════════════════════════════

[supervisord]
logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=%(ENV_HOME)s/mikkieworld/pids/supervisord.pid
nodaemon=false
silent=false
user=%(ENV_USER)s

; ─── Web interface (optioneel) ─────────────────────────────────────────
[inet_http_server]
port=127.0.0.1:9001
; username=mikkie
; password=geheim

[supervisorctl]
serverurl=http://127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ─── BRAIN — Centrale Orchestrator ────────────────────────────────────
[program:brain]
command=python3 %(ENV_HOME)s/mikkieworld/mikkie_brain.py loop
directory=%(ENV_HOME)s/mikkieworld
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stderr_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/brain_stderr.log
stdout_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/brain_stdout.log
stderr_logfile_maxbytes=5MB
stdout_logfile_maxbytes=5MB
environment=XAI_API_KEY="%(ENV_XAI_API_KEY)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",TELEGRAM_CHAT_ID="%(ENV_TELEGRAM_CHAT_ID)s",GUMROAD_API_TOKEN="%(ENV_GUMROAD_API_TOKEN)s"

; ─── TELEGRAM COMMANDER — Mobiele Controle ────────────────────────────
[program:commander]
command=python3 %(ENV_HOME)s/mikkieworld/mikkie_telegram_commander.py
directory=%(ENV_HOME)s/mikkieworld
autostart=true
autorestart=true
startretries=5
startsecs=3
stopwaitsecs=10
stderr_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/commander_stderr.log
stdout_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/commander_stdout.log
stderr_logfile_maxbytes=5MB
stdout_logfile_maxbytes=5MB
environment=XAI_API_KEY="%(ENV_XAI_API_KEY)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",TELEGRAM_CHAT_ID="%(ENV_TELEGRAM_CHAT_ID)s"

; ─── HEALER — Zelf-herstellend Systeem ────────────────────────────────
[program:healer]
command=python3 %(ENV_HOME)s/mikkieworld/mikkie_healer.py start
directory=%(ENV_HOME)s/mikkieworld
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=30
stderr_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/healer_stderr.log
stdout_logfile=%(ENV_HOME)s/MIKKIE_WORLD/LOGS/healer_stdout.log
stderr_logfile_maxbytes=5MB
stdout_logfile_maxbytes=5MB
environment=XAI_API_KEY="%(ENV_XAI_API_KEY)s",TELEGRAM_BOT_TOKEN="%(ENV_TELEGRAM_BOT_TOKEN)s",TELEGRAM_CHAT_ID="%(ENV_TELEGRAM_CHAT_ID)s"

; ─── Groepsbeheer ──────────────────────────────────────────────────────
[group:mikkie_daemons]
programs=brain,commander,healer
priority=999
